# demos/filemover/filemover.path
[Unit]
Description=Filemover watcher

[Path]
DirectoryNotEmpty=/var/tmp/in
#Unit=filemover.service

[Install]
WantedBy=multi-user.target
