# demos/filemover/filemover.service
[Unit]
Description=Filemover service

[Service]
ExecStart=/root/demos/filemover/moveit.sh
Type=oneshot
